# deploy/services/macos/com.skill.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.skill.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Skill.app/Contents/MacOS/skill-daemon</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/skill-daemon.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/skill-daemon.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>SKILL_DAEMON_ADDR</key>
    <string>127.0.0.1:18444</string>
  </dict>
</dict>
</plist>
